# jeannie-booking-bot/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

; Web process — gunicorn serves the Flask app with APScheduler (handles noon booking)
[program:gunicorn]
command=gunicorn --bind 0.0.0.0:%(ENV_PORT)s --workers 1 --timeout 120 app:app
directory=/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=5

; Worker process — external cron: refreshes token every 5 min, books daily at noon
; Runs AFTER gunicorn is up; no RAILWAY_PUBLIC_DOMAIN → defaults to http://127.0.0.1:5000
[program:cron]
command=python cron.py
directory=/app
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startretries=5
